# watering/sensors/systemd/water.timer
[Unit]
Description=test
Requires=water.service

[Timer]
Unit=water.service
OnUnitActiveSec=10min
OnBootSec=1min

[Install]
WantedBy=timers.target
