# watering/sensors/systemd/water.timer
[Unit]
Description=test

[Timer]
OnUnitActiveSec=10min
OnBootSec=1min

[Install]
WantedBy=water.target
